# bin/webpc-master.service
[Unit]
Description=WebPC Master Service
After=network.target
 
[Service]
#User=webpc
Type=simple
ExecStart=/opt/webpc/webpc master
KillMode=control-group
Restart=on-failure
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target